# scripts/agents/companion/com.atrophy.introspect.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.atrophy.introspect</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/williamlilley/.pyenv/versions/3.12.7/bin/python</string>
        <string>/Users/williamlilley/Projects/Claude Code Projects/Atrophy App Electron/scripts/agents/companion/introspect.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/williamlilley/Projects/Claude Code Projects/Atrophy App Electron</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>17</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/williamlilley/.atrophy/logs/introspect.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/williamlilley/.atrophy/logs/introspect.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/williamlilley/.pyenv/versions/3.12.7/bin</string>
        <key>AGENT</key>
        <string>companion</string>
    </dict>
</dict>
</plist>
